# association-mining.service
[Unit]
Description=Association Rule Mining Web Application
After=network.target

[Service]
Type=simple
User=duy
WorkingDirectory=/home/duy/KHDL_demo_main
Environment=PATH=/home/duy/KHDL_demo_main/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/duy/KHDL_demo_main
ExecStart=/home/duy/KHDL_demo_main/venv/bin/python app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
